# bin/lkp-setup-rootfs
#!/bin/sh
### BEGIN INIT INFO
# Provides:          nfsroot
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: NFSROOT setup
# Description:       NFSROOT setup
### END INIT INFO

umask 002

. $LKP_SRC/lib/bootstrap.sh

# initiation at boot stage; should be invoked once for
# each fresh boot.
boot_init()
{
	mount_tmpfs
	redirect_stdout_stderr

	setup_hostname
	setup_hosts

	announce_bootup

	add_lkp_user

	fixup_packages

	setup_network
	run_ntpdate

	mount_debugfs
	mount_rootfs

	export PATH=$PATH:$LKP_SRC/bin
	echo   PATH=$PATH
}

boot_init

echo $LKP_DEBUG_PREFIX $LKP_SRC/bin/run-lkp
$LKP_DEBUG_PREFIX $LKP_SRC/bin/run-lkp
