# bin/lkp-setup-rootfs
#!/bin/sh
### BEGIN INIT INFO
# Provides:          nfsroot
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: NFSROOT setup
# Description:       NFSROOT setup
### END INIT INFO

umask 002

export PATH=$PATH:$LKP_SRC/bin
echo   PATH=$PATH

. $LKP_SRC/lib/bootstrap.sh
. $LKP_SRC/lib/http.sh
. $LKP_SRC/lib/reboot.sh
. $LKP_SRC/lib/ucode.sh
. $LKP_SRC/lib/env.sh

# https://stackoverflow.com/questions/19100708/mongodb-mongorestore-failure-localefacet-s-create-c-locale-name-not-valid
# To fix:
# terminate called after throwing an instance of 'std::runtime_error'
#  what():  locale::facet::_S_create_c_locale name not valid
[ -z "$LC_ALL" ] && locale -a 2>/dev/null | grep -qw C.UTF-8 && export LC_ALL=C.UTF-8

is_docker || {
	# disable oom killing entirely for the task
	cat > /proc/$$/oom_score_adj <<EOF
-1000
EOF
}

boot_init || {
	echo "LKP: rebooting due to boot init failure"
	set_tbox_wtmp 'rebooting'
	start_reboot_watchdog
	reboot_tbox 2>/dev/null && exit 1
}

# when to not reboot is a policy issue, which can be refined later.
# here just serves as an example only.
# FIXME: need to consider more factors, such as cmdline, boot_parameter
job_does_not_need_reboot()
{
	grep -q "^do_not_reboot: 1$" $job && {
		echo "LKP: [do_not_reboot] start a new job without rebooting"
		return 0
	}

	grep -q "^do_not_reboot_for_same_kernel: 1$" $job &&
	is_same_kernel_and_rootfs && is_same_testcase && is_same_bp_memmap && {
		echo "LKP: [do_not_reboot_for_same_kernel] start a new job without rebooting"
		return 0
	}

	return 1
}

job_force_reboot()
{
	grep -q "^force_reboot: \+1$" $job
}

reboot_for_next_job()
{
	grep -q "^reboot_for_next_job: \+1$" $job
}

start_reboot_watchdog()
{
	[ -n "$boot_timeout" ] || return
	[ -f /var/lock/bmc-watchdog ] && return

	bmc-watchdog --set --start-after-set -i $boot_timeout
}

downgrade_ucode()
{
	local next_version=$(awk '/^ucode: /{print $2}' $job | tr -d "'")
	ucode_less_than_current $next_version
}

wget_install_cgz()
{
	local url
	for url
	do
		echo "wget $url"
		wget $url -O /tmp/tmp.cgz
		(
			cd /
			gzip -dc /tmp/tmp.cgz | cpio -id
		)
	done
	rm -f /tmp/tmp.cgz
}

install_proxy_yum()
{
	[ -d "/etc/yum/" ] || return

	echo "proxy=http://$SQUID_HOST:$SQUID_PORT" >> /etc/yum.conf
	if ls /etc/yum.repos.d/*.repo >/dev/null 2>&1; then
		sed -i 's/https:/http:/g' /etc/yum.repos.d/*.repo
	fi

	return 0
}

install_proxy_apt()
{
	[ -d "/etc/apt/" ] || return

	cat << EOF > /etc/apt/apt.conf.d/proxy
Acquire::http::proxy "http://$SQUID_HOST:$SQUID_PORT/";
Acquire::ftp::proxy "ftp://$SQUID_HOST:$SQUID_PORT/";
EOF

	for cfg in /etc/apt/sources.list /etc/apt/sources.list.d/*
	do
		[ -f "$cfg" ] || continue
		sed -i 's/https:/http:/g' $cfg
	done

	return 0
}

install_docker_proxy()
{
	[ -n "$DOCKER_PROXY_HOST" ] || return
	[ -n "$DOCKER_PROXY_PORT" ] || return

	mkdir -p /etc/systemd/system/docker.service.d
	cat << EOF > /etc/systemd/system/docker.service.d/http-proxy.conf
[Service]
Environment="HTTPS_PROXY=http://$DOCKER_PROXY_HOST:$DOCKER_PROXY_PORT/"
EOF

	curl http://$DOCKER_PROXY_HOST:$DOCKER_PROXY_PORT/ca.crt > \
		/etc/pki/ca-trust/source/anchors/docker_registry_proxy.crt || return
	update-ca-trust extract

	systemctl daemon-reload
	has_cmd docker && systemctl restart docker.service

	return 0
}

modify_gem_source()
{
       echo ':sources: ["http://rubygems.org"]' >> ~/.gemrc

       return 0
}

install_gitcache()
{
       [ -n "$gitcache_host" ] || return

       # the "https://git" prefix catches github/gitlab/gitee
       echo -e "[url \"http://$gitcache_host:${gitcache_port:-5000}/git\"]\n\tinsteadOf = https://git" >> /etc/gitconfig

       return 0
}

install_proxy()
{
	install_gitcache

	modify_gem_source

	install_docker_proxy

	[ -z "$SQUID_HOST" ] && return
	[ -z "$SQUID_PORT" ] && return

	install_proxy_apt || install_proxy_yum
}

# Save the result to
# "/srv/initrd/pkg/$os/$os_arch/$os_version/$benchmark.cgz",
# get the result through wget $initrd_pkg.
define_files_pkg_pack()
{
	(
		local benchmark=${1#pkg/}
		export benchmark=${benchmark%/PKGBUILD}
		$LKP_SRC/tests/makepkg
	)
}

# There may be $program-dev or $program-dev-dev in the
# "$LKP_SRC/distro/depends" directory,
# when the value passed by $define_files is "distro/depends/$program-dev",
# you don't know which script in "makepkg" or "pack-deps" to execute,
# so it needs to be handled sparately.
# Save the result to
# "/srv/initrd/deps/$os/$os_arch/$os_version/$deps_benchmark.cgz",
# get the result through wget $initrd_deps.
define_files_deps_pack()
{
	local deps_benchmark=${1##*/}
	local pkg_benchmark=${deps_benchmark%-dev}
	local suffix=${deps_benchmark##*-}

	[ $suffix == "dev" ] && [ -f "$LKP_SRC/pkg/$pkg_benchmark/PKGBUILD" ] && {
		(
			export benchmark=$pkg_benchmark
			$LKP_SRC/tests/makepkg
		)
		return 0
	}

	(
		export benchmark=$deps_benchmark
		export pack_arch=$os_arch
		$LKP_SRC/tests/pack-deps
	)
}

# $define_files may be like this:
# "pkg/$program/PKGBUILD distro/depends/$program distro/depends/$program-dev".
define_files_auto_pack()
{
	for path
	do
		[ -f "$LKP_SRC/$path" ] || continue

		[ ${path##*/} = PKGBUILD ] && {
			define_files_pkg_pack $path
		}
		[ ${path#distro/depends/} != $path ] && {
			define_files_deps_pack $path
		}
	done
}

install_proxy

while true; do
	set_tbox_wtmp 'running'
	if [ "$os_mount" != "initramfs" ]; then
		. $job_script define_files
		define_files_auto_pack $define_files
		wget_install_cgz $initrd_deps
		wget_install_cgz $initrd_pkg
		install_rpms
		install_deb
	fi

	echo $LKP_DEBUG_PREFIX $LKP_SRC/bin/run-lkp
	$LKP_DEBUG_PREFIX $LKP_SRC/bin/run-lkp $job

	is_docker && {
		echo "LKP: exiting"
		exit
	}
	
	if tbox_cant_kexec || job_force_reboot; then
		echo "LKP: rebooting"
		set_tbox_wtmp 'rebooting'
		# Avoid get stuck in yocto minimal rootfs.
		# LKP: rebooting
		# [    8.124081] LKP: rebooting
		# # (shell prompt)
		is_virt && [ -e '/proc/sysrq-trigger' ] && {
			sync
			echo s > /proc/sysrq-trigger
			echo b > /proc/sysrq-trigger
		}

		start_reboot_watchdog
		reboot_tbox 2>/dev/null && exit
	fi

	if [ -n "$LKP_DEBUG" ]; then
		# To debug what left from last run
		set -x
		ps faux 2>/dev/null || ps
		ls -l /tmp/
		ls -l /tmp/lkp
		lsmod
	fi

	reboot_for_next_job
	need_reboot=$?

	next_job
	rsync_rootfs
	[ "$need_reboot" -eq 0 ] && break
	job_does_not_need_reboot || break
done

if downgrade_ucode; then
	echo "reboot tbox to downgrade ucode"
	set_job_state 'downgrade_ucode'
	sleep 5
	reboot_tbox 2>/dev/null && exit
fi

set_tbox_wtmp 'kexec_to_next_job'
$LKP_SRC/bin/kexec-lkp ${pxe_user:-lkp} $NEXT_JOB
