# rootfs/addon/etc/systemd/system/rc-local.service.d/override.conf
[Service]
IgnoreSIGPIPE=no
TasksMax=infinity
